# crontab
15 06 * * * python /app/EDRAutomationTesting.py >> /var/log/cron.log 2>&1

